# etc/init/power-meter-webserver.conf
# ssh - OpenBSD Secure Shell server
#
# The OpenSSH server provides secure shell access to the system.

description	"Power Meter Express Webserver"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

expect fork

env APP_DIR="/home/tm/src/tfmalt/power-meter-webserver"
env APP_START="app.js"
env LOG_DIR="/var/log/power-meter"

pre-start script
    cd $APP_DIR
    pwd
end script

script
    exec forever --sourceDir $APP_DIR -a -l $LOG_DIR/express-access.log \
         --minUptime 5000 --spinSleepTime 2000 start $APP_START
end script

pre-stop script 
    exec forever stop $APP_START >> $LOG_DIR/express-access.log
end script
